# systembackup.service
[Unit]
Description=Incremental Backup

[Service]
CPUAffinity=0
Type=simple
ExecStart=/usr/bin/backup/systembackup.sh
