# systembackup.service
[Unit]
Description=Incremental Backup

[Service]
Type=simple
ExecStart=/home/owen/bin/backup/systembackup.sh
